# control/infolab-lights-restart.service
[Unit]
Description=Restart Infolab Lights Service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart infolab-lights.service